# src/src/pcp/openlmi-pcp.cron
#! /bin/sh

openlmi-pcp-generate localhost >/dev/null 2>&1
