# src/src/service-dbus/test/failing.service
[Unit]
Description=unit %n

[Service]
Type=oneshot
ExecStart=-/usr/bin/false
ExecStop=-/usr/bin/false
